# config/launchd/com.astock.trade.health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.astock_trading.trade.health</string>
  <key>WorkingDirectory</key>
  <string>/path/to/a-stock-trading</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>ASTOCK_DATABASE_URL</key>
    <string>mysql+pymysql://user:password@host:3306/a_stock_trading?charset=utf8mb4</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/path/to/a-stock-trading/bin/trade</string>
    <string>health</string>
    <string>--json</string>
  </array>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StandardOutPath</key>
  <string>/path/to/a-stock-trading/logs/launchd/health.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/a-stock-trading/logs/launchd/health.err.log</string>
</dict>
</plist>
